# devops/ansible/nodeexporter.service
[Unit]
Description=NodeExporter

[Service]
TimeoutStartSec=0
ExecStart=/var/nodeexporter/node_exporter

[Install]
WantedBy=multi-user.target
